# launchd/com.mugi-claw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mugi-claw</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>--env-file=.env</string>
        <string>dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kagadminmac/project/mugi-claw</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/kagadminmac/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/kagadminmac/.mugi-claw/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kagadminmac/.mugi-claw/launchd-stderr.log</string>
</dict>
</plist>
